# systemd/BikeCon-web.service
# /etc/systemd/system/BikeCon-web.service
[Unit]
Description=BikeCon Web Dashboard
After=BikeCon-mixer.service network.target

[Service]
Type=simple
User=root
Environment="BIKEYCON_HOME=/opt/BikeCon"
Environment="BIKEYCON_CONFIG=/etc/BikeCon"
WorkingDirectory=$BIKEYCON_HOME
ExecStart=$BIKEYCON_HOME/c2lite/bin/python3 -m uvicorn webapp:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
