# systemd/BikeCon-joycon.service
# /etc/systemd/system/BikeCon-joycon.service
[Unit]
Description=BikeCon Joy-Con Driver
Requires=BikeCon-mixer.service
After=BikeCon-mixer.service bluetooth.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/BikeCon
ExecStart=/opt/BikeCon/venv/bin/python3 joycon_service.py
Restart=always
RestartSec=5
ExecStopPost=-/bin/sh -c 'for mac in $(/usr/bin/bluetoothctl devices Connected | grep -i "Joy-Con" | awk "{print $$2}"); do /usr/bin/bluetoothctl disconnect $mac; done'
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
