# systemd/c2lite-bike.service
# /etc/systemd/system/c2lite-bike.service
[Unit]
Description=C2Lite Bike Bluetooth Driver
Requires=c2lite-mixer.service
After=c2lite-mixer.service bluetooth.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/zhou/c2lite_controller
ExecStart=/home/zhou/c2lite_controller/c2lite/bin/python3 bike_service.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
KillSignal=SIGTERM
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target